# launchd/com.nanoclaw.venice-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.venice-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{NPX_PATH}}</string>
        <string>tsx</string>
        <string>proxy/venice-proxy.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{NODE_BIN_DIR}}:{{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>VENICE_API_KEY</key>
        <string>{{VENICE_API_KEY}}</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/venice-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/venice-proxy.error.log</string>
</dict>
</plist>
